# deploy/termigate.service
[Unit]
Description=termigate — Remote terminal manager
After=network.target

[Service]
Type=exec
User=ben
Group=ben
WorkingDirectory=/opt/termigate
Environment=HOME=/home/ben
Environment=PORT=4000
Environment=PHX_HOST=localhost
Environment=PHX_SERVER=true
Environment=SECRET_KEY_BASE=generate-a-secret-key-base-here
# Uncomment for remote access:
# Environment=PHX_BIND=0.0.0.0
# Environment=TERMIGATE_AUTH_TOKEN=your-secure-token
ExecStart=/opt/termigate/bin/termigate start
ExecStop=/opt/termigate/bin/termigate stop
Restart=on-failure
RestartSec=5
# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/tmp/termigate /home/ben/.config/termigate
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
