# deploy/tmux-rm.service
[Unit]
Description=tmux-rm — Remote terminal manager
After=network.target

[Service]
Type=exec
User=ben
Group=ben
WorkingDirectory=/opt/tmux_rm
Environment=HOME=/home/ben
Environment=PORT=4000
Environment=PHX_HOST=localhost
Environment=PHX_SERVER=true
Environment=SECRET_KEY_BASE=generate-a-secret-key-base-here
# Uncomment for remote access:
# Environment=PHX_BIND=0.0.0.0
# Environment=RCA_AUTH_TOKEN=your-secure-token
ExecStart=/opt/tmux_rm/bin/tmux_rm start
ExecStop=/opt/tmux_rm/bin/tmux_rm stop
Restart=on-failure
RestartSec=5
# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/tmp/tmux-rm /home/ben/.config/tmux_rm
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
